# system_programing/deamon/myservice.service
[Unit]
Description=Моя простая служба

[Service]
ExecStart=/opt/myservice.sh
Restart=always

[Install]
WantedBy=multi-user.target
